# KoboServer/src/etc/udev/rules.d/98-koboserver.rules
KERNEL=="eth*", ACTION=="add", RUN+="/usr/local/koboserver/udev_server.sh"
KERNEL=="wlan*", ACTION=="add", RUN+="/usr/local/koboserver/udev_server.sh"

